# recipes-support/startup-container/files/start-container@.service
[Unit]
Description=Run docker container %I
Requires=docker.service
After=docker.service

[Service]
Type=forking
ExecStart=/usr/bin/docker start %i
ExecStop=/usr/bin/docker stop -t 5 %i

[Install]
WantedBy=multi-user.target
